# deploy/systemd/reames-agent.service
[Unit]
Description=Reames Agent Serve
After=network.target

[Service]
Type=simple
User=reames
WorkingDirectory=/opt/reames-agent
EnvironmentFile=-/opt/reames-agent/.env
ExecStart=/opt/reames-agent/bin/reames-agent serve --addr 127.0.0.1:8787
Restart=always
RestartSec=5
Environment=REAMES_AGENT_HOME=/opt/reames-agent/data
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
